# RIScraper/crontab
* */3 * * * python /RIScraper/Scraper.py